# deploy/camera-brain-chat.service
[Unit]
Description=Camera Brain Chat Web Interface
After=network.target llama-server.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/camera-brain
Environment="PATH=/home/camera-brain/chat-env/bin:/usr/bin:/bin"
Environment="DB_PATH=/home/camera-brain/camera-brain.db"
Environment="LLAMA_SERVER_URL=http://localhost:8889"
Environment="MODEL_PATH=/home/camera-brain/models/LFM2.5-1.2B-Instruct.Q4_K_M.gguf"
Environment="PORT=8080"
ExecStart=/home/camera-brain/chat-env/bin/python /home/camera-brain/chat/app.py
Restart=on-failure
RestartSec=5

# Security hardening
PrivateTmp=true

[Install]
WantedBy=multi-user.target
