# deploy/camera-brain-gateway.service
[Unit]
Description=Camera Brain Gateway Service
After=network.target postgresql.service nats.service

[Service]
Type=simple
User=caimlas
Group=caimlas
Environment="WORKER_MODEL_PATH=/home/camera-brain/models/yolov5s_int8.rknn"
Environment="DB_PASSWORD=camera_brain_password_change_me"
Environment="DB_HOST=localhost"
Environment="DB_PORT=5432"
Environment="DB_USER=camera_brain"
Environment="DB_NAME=camera_brain"
Environment="NATS_URL=nats://localhost:4222"
Environment="SERVICE_PORT=8080"
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"
WorkingDirectory=/home/camera-brain
ExecStart=/home/camera-brain/gateway/gateway
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
